# plugins/daemon-manager/data/gxde-daemon-manager.service
[Unit]
Description=GXDE Daemon Manager (plugin registry and supervisor)
Documentation=https://gitee.com/GXDE-OS/gxde-daemon
# Bind to the graphical session: started on login, stopped on logout.
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=dbus
BusName=top.gxde.daemon.manager
ExecStart=/usr/libexec/gxde-daemon/daemon-manager/gxde-daemon-manager
# The manager restarts crashed plugins itself; this restarts the manager too.
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
