# misc/autostart/gxde-daemon-manager.desktop
[Desktop Entry]
Type=Application
Name=GXDE Daemon Manager
Name[zh_CN]=GXDE 守护进程管理器
Name[zh_TW]=GXDE 守護行程管理器
Comment=Start the GXDE daemon plugin registry and supervisor on login
Comment[zh_CN]=登录时启动 GXDE 守护进程插件注册与监管服务
Comment[zh_TW]=登入時啟動 GXDE 守護行程外掛註冊與監管服務
# The GXDE session does not reach graphical-session.target, so the systemd
# user unit is triggered from XDG autostart instead. Going through systemctl
# (rather than exec-ing the binary directly) keeps the unit's Restart=on-failure
# supervision in effect.
Exec=systemctl --user start gxde-daemon-manager.service
Terminal=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
OnlyShowIn=Deepin;GXDE;
